# data/ps_coor.service
[Unit]
Description=ps_coor
Before=nginx.service
BindsTo=nginx.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ps_coor_starter.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
WantedBy=nginx.service
